# service/ipv6.service
[Unit]
After=network.service

[Service]
ExecStart=/usr/bin/ipv6.sh

[Install]
WantedBy=default.target
